# units/replace_private_ipv4.service
[Unit]
Description=Replace private_ipv4

[Service]
Type=oneshot

RemainAfterExit=yes

ExecStart=/etc/private_ipv4_replacer.sh
